# templates/mtproxy.service
[Unit]
Description=Telegram MTProxy (Managed Fake TLS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/mtproxy/mtproxy.env
ExecStart=/usr/local/bin/mtproxy-managed-run.sh
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target

